# pkg.FreeBSD/2.2.8/pkg/PLIST
etc/rc.d/bkpupsd.sh.sample
etc/upstab
sbin/bkpupsd
libexec/bkpupsd/battery.shutdown
libexec/bkpupsd/line.lost
libexec/bkpupsd/line.recover
libexec/bkpupsd/line.shutdown
@dirrm libexec/bkpupsd
